# scripts/back_end.service
[Unit]
Description=Spark service

[Service]
ExecStart=/home/administrator/configuratie/scripts/back_end.sh

[Install]
WantedBy=multi-user.target
